# youtube-bot.service
[Unit]
Description=YouTube to MP3 Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/youtube-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/youtube-agent/venv/bin/python3 youtube_telegram_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 